# s6/sleep/run
#!/bin/sh

echo "$KEY should be VALUE"
echo "stderr" >&2

exec sleep 1000
