# s6/sleep/run
#!/bin/sh

echo "what"
echo "do"

exec sleep 1000
